# repos/trident_innate-os_ros2_robot_operating_system/systemd/zenoh-router.service
[Unit]
Description=Zenoh Router for ROS2
After=network-online.target
Wants=network-online.target

[Service]
User=jetson1
Group=jetson1

ExecStart=/home/jetson1/innate-os/dds/start_zenoh_router.zsh

# Restart the service if it fails
Restart=on-failure
RestartSec=5

# Standard output goes to journald
StandardOutput=journal
StandardError=journal

# Give it time to bind to port
TimeoutStartSec=10

# Don't wait too long during shutdown
TimeoutStopSec=10
# Kill main process with SIGTERM, then SIGKILL remaining children
KillMode=mixed

[Install]
WantedBy=multi-user.target